# docker/app/shared/supervisor/webserver.conf
[program:nailsinc]
command=uwsgi --ini /var/www/conf/uwsgi.ini
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/host/logs/uwsgi.log
redirect_stderr=True
